# deploy/standard/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/foxnas/supervisord.log
pidfile=/var/run/supervisord.pid

[program:mysql]
command=/usr/sbin/mysqld --user=mysql --datadir=/var/lib/mysql --bind-address=0.0.0.0
priority=10
autostart=true
autorestart=true
stdout_logfile=/var/log/foxnas/mysql.log
stderr_logfile=/var/log/foxnas/mysql.err.log

[program:redis]
command=/usr/bin/redis-server --appendonly yes --dir /data/redis --bind 0.0.0.0
priority=20
autostart=true
autorestart=true
stdout_logfile=/var/log/foxnas/redis.log
stderr_logfile=/var/log/foxnas/redis.err.log

[program:rabbitmq]
command=/usr/sbin/rabbitmq-server
priority=30
autostart=true
autorestart=true
stdout_logfile=/var/log/foxnas/rabbitmq.log
stderr_logfile=/var/log/foxnas/rabbitmq.err.log

[program:foxnas]
command=/opt/foxnas/start-foxnas.sh
priority=40
autostart=true
autorestart=true
stdout_logfile=/var/log/foxnas/backend.log
stderr_logfile=/var/log/foxnas/backend.err.log

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
priority=50
autostart=true
autorestart=true
stdout_logfile=/var/log/foxnas/nginx.log
stderr_logfile=/var/log/foxnas/nginx.err.log

